# crontab
10 11 * * *   ~/atlas-dyndns/topology-update.sh >> ~/log/topology-update.log
